# roles/cloudflare-ddns/templates/cf_ddns.service.j2
[Unit]
Description=Update DNS entry for this host to current IP

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 {{ cf_ddns_path }}pyflare.py
