# roles/cloudflare-ddns/templates/cf_ddns.service
[Unit]
Description=Update DNS entry for this host to current IP

[Service]
Type=oneshot
ExecStart=python {{ cf_ddns_path }}/pyflare.py
